# etc/inittab
id:1:initdefault:
fs::bootwait:/bin/readfstab.lua
1:1:respawn:/bin/getty.lua tty1
